# supervisor_glare.conf
[program:glare]
command=gunicorn -k gevent -w 2 -b 0.0.0.0:80 main:app
directory=/glare/
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
